# contrib/alpine/rc.local
#!/bin/sh
zcat /etc/keymap/se.bmap.gz | loadkmap
