# contrib/alpine/rc.local
#!/bin/sh
zcat /etc/keymap/sv-latin1.bmap.gz | loadkmap
