# contrib/void/rc.local
#!/bin/sh
# This script is executed once at boot, at the very end,
# right before starting the TTYs.
#
# Each line is run in sequence.
#
# NOTE: Remember to set the executable bit!
